# build/anshuman.me.conf
[program:anshuman.me]
command=/opt/anshuman.me/my-portfolio/deploy/start_gunicorn
directory=/opt/anshuman.me/my-portfolio/
autostart=true
autorestart=true